# scripts/com.jvdh.tracks-job-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <!-- Template only. Use scripts/setup_machine.sh and scripts/install_scheduler.sh. -->
  <dict>
    <key>Label</key>
    <string>com.jvdh.jobsearch.scheduler.__JOBSEARCH_SCHEDULER_INSTANCE__</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__JOBSEARCH_ROOT__/scripts/run_scheduled_jobs.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>__HOME__</string>
      <key>PATH</key>
      <string>__PATH__</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__JOBSEARCH_ROOT__</string>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>__JOBSEARCH_ROOT__/logs/scheduler.out</string>

    <key>StandardErrorPath</key>
    <string>__JOBSEARCH_ROOT__/logs/scheduler.err</string>
  </dict>
</plist>
